# .config/systemd/user/awww-daemon.service
[Unit]
Description=AWWW wallpaper daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/awww-daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
