# .config/systemd/user/battery-warning.service
[Unit]
Description=Battery warning check

[Service]
Type=oneshot
ExecStart=%h/.config/i3/battery-warning.sh
