# .config/systemd/user/udiskie.service
[Unit]
Description=Udiskie automount daemon
Documentation=https://github.com/coldfix/udiskie
After=default.target

[Service]
Type=simple
ExecStart=/usr/bin/udiskie -n
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
