# .config/systemd/user/dropbox.service.d/override.conf
[Service]
ExecStart=
ExecStart=/home/harry/.dropbox-dist/dropbox-lnx.x86_64-252.4.3485/dropboxd
KillMode=mixed
Restart=on-failure
